# com.example.bhyve-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd template for continuous B-hyve monitoring.
     1. Replace USERNAME and the two paths below.
     2. cp to ~/Library/LaunchAgents/com.example.bhyve-monitor.plist
     3. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.example.bhyve-monitor.plist
     Stop/remove: launchctl bootout gui/$(id -u) <plist> -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.bhyve-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/GitHub/bhyve/.venv/bin/python</string>
        <string>/Users/USERNAME/GitHub/bhyve/monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/GitHub/bhyve</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/bhyve-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/bhyve-monitor.err.log</string>
</dict>
</plist>
